# docs/launchd/stock-reminder-open-sync.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>local.stock-reminder.open-sync</string>

  <key>WorkingDirectory</key>
  <string>/Users/liuxiang/Programme/Self_Play/Stock_Reminder</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/liuxiang/Programme/Self_Play/Stock_Reminder/scripts/run_daily_sync.sh</string>
    <string>open</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>DATABASE_PATH</key>
    <string>/Users/liuxiang/Programme/Self_Play/Stock_Reminder/backend/data/stock_reminder.sqlite3</string>
    <key>LOG_DIR</key>
    <string>/Users/liuxiang/Programme/Self_Play/Stock_Reminder/logs</string>
  </dict>

  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Weekday</key>
      <integer>1</integer>
      <key>Hour</key>
      <integer>9</integer>
      <key>Minute</key>
      <integer>40</integer>
    </dict>
    <dict>
      <key>Weekday</key>
      <integer>2</integer>
      <key>Hour</key>
      <integer>9</integer>
      <key>Minute</key>
      <integer>40</integer>
    </dict>
    <dict>
      <key>Weekday</key>
      <integer>3</integer>
      <key>Hour</key>
      <integer>9</integer>
      <key>Minute</key>
      <integer>40</integer>
    </dict>
    <dict>
      <key>Weekday</key>
      <integer>4</integer>
      <key>Hour</key>
      <integer>9</integer>
      <key>Minute</key>
      <integer>40</integer>
    </dict>
    <dict>
      <key>Weekday</key>
      <integer>5</integer>
      <key>Hour</key>
      <integer>9</integer>
      <key>Minute</key>
      <integer>40</integer>
    </dict>
  </array>

  <key>StandardOutPath</key>
  <string>/Users/liuxiang/Programme/Self_Play/Stock_Reminder/logs/launchd-open-sync.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/liuxiang/Programme/Self_Play/Stock_Reminder/logs/launchd-open-sync.err.log</string>
</dict>
</plist>
